# virus_hog/add_to_system.service
[Unit]
Description=Script Daemon

[Service]
Type=simple
ExecStart=etc/systemd/system/run.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
